# scripts/bizobs.service
[Unit]
Description=Business Observability Forge Server
Documentation=https://github.com/lawrobar90/Dynatrace-Business-Observability-Forge
After=network.target docker.service
Wants=network-online.target docker.service
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/Dynatrace-Business-Observability-Forge
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=HOME=/home/ec2-user
Environment=PATH=/home/ec2-user/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/home/ec2-user/Dynatrace-Business-Observability-Forge/.env
ExecStartPre=/bin/bash /home/ec2-user/Dynatrace-Business-Observability-Forge/scripts/self-heal.sh
ExecStart=/home/ec2-user/.nvm/versions/node/v24.14.0/bin/node --require ./otel.cjs server.js
Restart=always
RestartSec=5
SyslogIdentifier=bizobs-server
StandardOutput=journal
StandardError=journal
TimeoutStartSec=60
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/ec2-user/Dynatrace-Business-Observability-Forge
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
